# core-api/deploy/systemd/bodeul-core-api-preview.service.example
[Unit]
Description=BoDeul Core API Preview
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bodeul-api
Group=bodeul-api
WorkingDirectory=/opt/bodeul-core-api
EnvironmentFile=/etc/bodeul-core-api/preview.env
Environment=SPRING_PROFILES_ACTIVE=preview
ExecStart=/usr/bin/java -jar /opt/bodeul-core-api/current.jar
Restart=on-failure
RestartSec=5
SuccessExitStatus=143
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
